# perfSONAR_PS-LookupService/scripts/lookup_service
#!/bin/bash
#
# Init file for perfSONAR Lookup Service
#
# chkconfig: 2345 60 20
# description: perfSONAR Lookup Service
#

# get local lookup service configuration
if [ -f /etc/sysconfig/lookup_service ];then
        . /etc/sysconfig/lookup_service
fi

# defaults, if not specified in local configuration
[ "$RESET_DATABASE" = "" ] && RESET_DATABASE=NO
[ "$CLEAN_DATABASE" = "" ] && CLEAN_DATABASE=YES

PREFIX=/opt/perfsonar_ps/lookup_service
BINDIR=${PREFIX}/bin
CONFDIR=${PREFIX}/etc
RUNDIR=/var/run

CONFFILE=${CONFDIR}/daemon.conf
LOGGERFILE=${CONFDIR}/daemon_logger.conf
PIDFILE=lookup_service.pid

DATABASE_DIRECTORY=/var/lib/perfsonar/lookup_service/xmldb

USER=perfsonar
GROUP=perfsonar

FULL_PIDFILE=${RUNDIR}/lookup_service.pid

PERFSONAR="${BINDIR}/daemon.pl --config=${CONFFILE} --pidfile=${PIDFILE} --piddir=${RUNDIR} --logger=${LOGGERFILE} --user=${USER} --group=${GROUP}"

if [ "$CLEAN_DATABASE" == "NO" ]; then
	PERFSONAR="${PERFSONAR} --noclean"
fi

ERROR=0
ARGV="$@"
if [ "x$ARGV" = "x" ] ; then 
    ARGS="help"
fi

for ARG in $@ $ARGS
do
    # check for pidfile
    if [ -f $FULL_PIDFILE ] ; then
        PID=`cat $FULL_PIDFILE`
        if [ "x$PID" != "x" ] && kill -0 $PID 2>/dev/null ; then
            STATUS="perfSONAR Lookup Service (pid $PID) running"
            RUNNING=1
        else
            STATUS="perfSONAR Lookup Service (pid $PID?) not running"
            RUNNING=0
        fi
    else
        STATUS="perfSONAR Lookup Service (no pid file) not running"
        RUNNING=0
    fi 

    case $ARG in
    start)
	echo $PERFSONAR

	# Delete the database contents on startup
        if [ "$RESET_DATABASE" == "YES" ]; then
	    rm -rf $DATABASE_DIRECTORY/*
        fi

	if $PERFSONAR ; then
	    echo "$0 $ARG: perfSONAR Lookup Service started"
            touch /var/lock/subsys/lookup_service
	else
	    echo "$0 $ARG: perfSONAR Lookup Service could not be started"
	    ERROR=3
	fi
	;;
    stop)
	if [ $RUNNING -eq 0 ]; then
	    echo "$0 $ARG: $STATUS"
	    continue
	fi
	if kill $PID ; then
	    echo "$0 $ARG: perfSONAR Lookup Service stopped"
            rm /var/lock/subsys/lookup_service
	else
	    echo "$0 $ARG: perfSONAR Lookup Service could not be stopped"
	    ERROR=4
	fi
	;;
    restart)
    	$0 stop; echo "waiting..."; sleep 10; $0 start;
	;;
    *)
	echo "usage: $0 (start|stop|restart|help)"
	cat <<EOF

start      - start perfSONAR Lookup Service
stop       - stop perfSONAR Lookup Service
restart    - restart perfSONAR Lookup Service if running by sending a SIGHUP or start if 
             not running
help       - this screen

EOF
	ERROR=2
    ;;

    esac

done

exit $ERROR
